# etc/udev/rules.d/90-zt-bridge.rules
SUBSYSTEM=="net", RUN+="/usr/local/bin/zt-bridge.sh"
